# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/logs/supervisord.pid

[program:answer_service]
command=python -m uvicorn answer_service.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/answer_service.err.log
stdout_logfile=/app/logs/answer_service.out.log

[program:fusion_service]
command=python -m uvicorn fusion_service.main:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/fusion_service.err.log
stdout_logfile=/app/logs/fusion_service.out.log

[program:emotion_service]
command=python -m uvicorn emotion_service.main:app --host 0.0.0.0 --port 8002
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/emotion_service.err.log
stdout_logfile=/app/logs/emotion_service.out.log

[program:insight_service]
command=python -m uvicorn insight_service.main:app --host 0.0.0.0 --port 8003
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/insight_service.err.log
stdout_logfile=/app/logs/insight_service.out.log

[program:engagement_service]
command=python -m uvicorn engagement_service.main:app --host 0.0.0.0 --port 8004
directory=/app
autostart=true
autorestart=true
startretries=5
stderr_logfile=/app/logs/engagement_service.err.log
stdout_logfile=/app/logs/engagement_service.out.log

[program:streamlit]
command=python -m streamlit run demo_app.py --server.port=8501 --server.address=0.0.0.0 --server.headless=true
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/app/logs/streamlit.err.log
stdout_logfile=/app/logs/streamlit.out.log
